# rootfs/config/inittab
::sysinit:/sbin/openrc sysinit
::sysinit:/sbin/openrc boot
::wait:/sbin/openrc default
::shutdown:/opt/orb/vinit-stop
::shutdown:/sbin/openrc shutdown

::sysinit:/opt/orb/vinit-early
::wait:/opt/orb/vinit-late
