# rootfs/config/inittab
::sysinit:/sbin/openrc sysinit
::sysinit:/sbin/openrc boot
::wait:/sbin/openrc default
::shutdown:/sbin/openrc shutdown

::sysinit:/opt/vc/vinit-early
::wait:/opt/vc/vinit-late
